# packaging/debian/debian/openvstorage-health-check.cron.d
# /etc/cron.d/openvstorage-health-check: crontab entries for the openvstorage-health-check package

* *   * * *  root  /usr/bin/ovs healthcheck silent